# crontab
#Minute	Hours	Day	Month	Weekday	
*	0/2	*	*	*	/bin/run.sh
# Don't edit /bin/run.sh , timings can be edited just fine though.
